# files/services/cause-violation.service
[Unit]
Description=Cause SELinux violation

[Service]
Type=simple
ExecStart=/opt/tutorial/bin/cause-violation
